# etc/init.d/omninames.in
#!/bin/bash
#
# omninames:    Starts the omniNames server
#
# Author:       Duncan Grisby <duncan@tideway.com>
#
# chkconfig: 345 98 02
# description: omniORB name service
# processname: omniNames
# config: /etc/omniORB.cfg

# Source function library.
[ -f /etc/rc.d/init.d/functions ] || exit 0
. /etc/rc.d/init.d/functions

RETVAL=0

start() {
        HOSTNAME=`hostname`
	if [ ! -f /var/omninames/omninames-$HOSTNAME.log ]; then
	    OMNINAMES_START="-start"
	fi
	echo -n $"Starting omniNames: "
        daemon @prefix@/bin/omniNames-daemon $OMNINAMES_START -errlog /var/omninames/errlog
        RETVAL=$?
        echo
        [ $RETVAL -eq 0 ] && touch /var/lock/subsys/omninames
	return $RETVAL
}

stop() {
        echo -n $"Stopping omniNames: "
        killproc omniNames
        RETVAL=$?
        [ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/omninames
        echo
	return $RETVAL
}

# See how we were called.
case "$1" in
  start)
	start
	;;
  stop)
	stop
	;;
  restart|reload)
	stop
	start
	;;
  condrestart)
	if [ -f /var/lock/subsys/omninames ]; then
	    stop
	    start
	fi
	;;
  *)
	echo $"Usage: $0 {start|stop|restart|reload|condrestart}"
	exit 1
esac

exit $RETVAL
